# docker_file/etc/profile.d/04-llvm.sh
#!/bin/bash
# set environmental variable for intel-compiler
export ONEAPI_ROOT=/opt/intel/oneapi
#!/bin/bash
export PATH=/usr/local/bin/llvm:${PATH}